# dbt/crontab
*/5 * * * * root cd /app && dbt run >> /var/log/cron/dbt.log 2>&1
 